# tools/systemd/homelab-deploy.timer
[Unit]
Description=Run the homelab GitOps reconcile on a schedule

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
Persistent=true

[Install]
WantedBy=timers.target
